# archlinux/services/htop.service
[Unit]
Description=htop on tty11

[Service]
Type=simple
ExecStart=/usr/bin/htop
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty11
